# docker/supervisor-logstash.conf
[program:logstash]
command=/opt/logstash/start-logstash.sh
stdout_logfile=/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s-stderr.log
