# docker/supervisor-cron.conf
[program:cron]
command=/usr/sbin/cron
stdout_logfile=/var/log/supervisor/%(program_name)s-stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s-stderr.log
